# infraestrutura/systemd/rapidleech.service
[Unit]
Description=RED Rapidleech
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rapidleech
Environment=RAPIDLEECH_HOST=127.0.0.1
Environment=RAPIDLEECH_PORT=2581
EnvironmentFile=-/etc/red-rapidleech.env
ExecStart=/bin/sh -lc 'exec /usr/bin/php -S "${RAPIDLEECH_HOST}:${RAPIDLEECH_PORT}" -t .'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
