# infraestrutura/systemd/red-seb-webhook.service
[Unit]
Description=SEB Session Webhook to WhatsApp
After=network-online.target red-openclaw.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/red-seb-monitor
Environment=SEB_WEBHOOK_PORT=2590
Environment=RED_SEB_PUBLIC_URL=http://redsystems.ddns.net:2580
EnvironmentFile=-/etc/red-seb-webhook.env
EnvironmentFile=-/etc/red-openclaw.env
ExecStart=/usr/bin/node /opt/red-seb-monitor/webhook-whatsapp.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
